# supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/tmp/supervisord.log
pidfile=/tmp/supervisord.pid

[program:postgres]
command=/usr/lib/postgresql/15/bin/pg_ctl -D /var/lib/postgresql/15/main -l /tmp/postgresql.log start -w
autostart=false
autorestart=false

[program:redis]
command=redis-server --appendonly yes --daemonize no
autostart=true
autorestart=true

[program:traefik]
command=/usr/local/bin/traefik --configFile=/etc/traefik/traefik.yml
autostart=true
autorestart=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0

[program:auth-service]
command=npx tsx src/index.ts
directory=/app/services/auth-service
autostart=true
autorestart=true
environment=PORT="3010",JWT_SECRET="dev-secret-change-in-prod"
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0

[program:workspace-api]
command=npx tsx src/index.ts
directory=/app/services/workspace-api
autostart=true
autorestart=true
environment=PORT="3001",DATABASE_URL="postgresql://workspace:workspace_dev@localhost:5432/workspace",REDIS_URL="redis://localhost:6379",JWT_SECRET="dev-secret-change-in-prod",CLIENT_URL="http://localhost"
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0

[program:note-api]
command=npx tsx src/index.ts
directory=/app/services/note-api
autostart=true
autorestart=true
environment=PORT="3000",DATABASE_URL="postgresql://workspace:workspace_dev@localhost:5432/note",WORKSPACE_DATABASE_URL="postgresql://workspace:workspace_dev@localhost:5432/workspace"
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0

[program:kanban-api]
command=npx tsx src/index.ts
directory=/app/services/kanban-api
autostart=true
autorestart=true
environment=PORT="8000",KANBAN_DATABASE_URL="postgresql://workspace:workspace_dev@localhost:5432/kanban",DB_HOST="localhost",DB_PORT="5432",DB_USER="workspace",DB_PASSWORD="workspace_dev",DB_NAME="workspace"
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0

[program:web]
command=npm run preview -- --port 5173 --host
directory=/app/apps/workspace-web
autostart=true
autorestart=true
environment=VITE_API_URL="/ws",VITE_WS_URL="/ws"
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
